# modules/systemd/user/joystickwake@.service
[Unit]
Description=Joystickwake at %I

[Service]
Type=simple
ExecStart=/usr/bin/nj joystickwake --loglevel info
