# modules/systemd/user/mywallmovie@.service.in
[Unit]
Description=WallMovie at %I

[Service]
Type=simple
ExecStart=/usr/bin/sh @startups@/mywallmovie.sh
Restart=on-failure
RestartSec=1s
SyslogIdentifier=mywallmovie@%I
