# modules/systemd/user/timerhour.timer
[Unit]
Description=Timer hourly timer

[Timer]
OnCalendar=hourly
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
