# modules/systemd/user/timerday.timer
[Unit]
Description=Timer daily timer

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
